# docker/cron.d/download-data
50 10 * * * www-data cd /var/www && php artisan download-expedia-data content 4 
50 10 * * * www-data cd /var/www && php artisan download-giata-data 
