# docker/cron.d/download-data
0 1 * * * www-data cd /var/www && php artisan test-command
