# docker/cron.d/download-data
55 11 * * * www-data cd /var/www && php artisan download-expedia-data content 4 
55 11 * * * www-data cd /var/www && php artisan download-giata-data 
